# scripts/pgsql/crontab
0 5 * * * /var/lib/pgsql/vacuum_db.sh
10 5 * * * /var/lib/pgsql/pgsql_reindex.sh

